# .docker/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv sh

/usr/bin/caddy run --environ --config /etc/caddy/Caddyfile
