# .docker/rootfs/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv sh

exec /worker.sh
