# .docker/rootfs/etc/s6-overlay/s6-rc.d/php/run
#!/command/with-contenv sh

mkdir -p /run/php/
exec /usr/sbin/php-fpm8.3 --nodaemonize
